# systemv/rad
#!/bin/bash
# chkconfig: 2345 18 82

### BEGIN INIT INFO
# Provides: rad
# Required-Start:    $remote_fs $syslog $network
# Required-Stop:     $remote_fs $syslog $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Router Announcement 
# Description: Router Announcement Daemon 
### END INIT INFO

# source function library
. /lib/lsb/init-functions

RETVAL=0
PROG="rad"
cmd=/usr/sbin/rad
LOCK_FILE=/var/lock/rad
WATCH_OPTS="--loglevel info"


case "$1" in
  start)
	log_daemon_msg "Starting $PROG" "$PROG"
	mkdir -p /var/run/ratools
	mkdir -p /var/run/rad
	#start_daemon $cmd $WATCH_OPTS &
	/usr/sbin/rad WATCH_OPTS >/var/log/rad.log &
	sleep 1
	if [ -f /etc/rad.conf ]; then
	     echo "Configuring $PROG"
	     /usr/bin/ractl < /etc/rad.conf
	fi
	RETVAL=$?
	[ $RETVAL -eq 0 ] && touch $LOCK_FILE
	log_end_msg $RETVAL
	;;
  stop)
	log_daemon_msg "Shutting down $PROG" "$PROG"
	killall -w -s SIGINT $cmd
	RETVAL=$?
	[ $RETVAL -eq 0 ] && rm -f $LOCK_FILE
	log_end_msg $RETVAL
	;;
  restart|reload|force-reload)
	$0 stop
	$0 start
	RETVAL=$?
	;;
  condrestart|try-restart)
	if [ -f $LOCK_FILE ]; then
		$0 stop
		$0 start
	fi
	RETVAL=$?
	;;
  status)
	status_of_proc -p /run/$PROG.pid $cmd
	RETVAL=$?
	;;
  *)
	echo $"Usage: $0 {start|stop|restart|reload|force-reload|condrestart|try-restart|status}"
	exit 2
esac

exit $RETVAL
